# rocon/devel/etc/catkin/profile.d/15.rospair.bash
# generated from rocon_python_comms/env-hooks/15.rospair.bash.em

. "/home/sclab_robot/turtlebot_ws/rocon/src/rocon_tools/rocon_python_comms/completion/rospair.bash"
